# recipes-beacon/beacon-eeprom-services/dump_eeprom.service
[Unit]
Description=Dump EEPROM contents

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/dump-eeprom.sh

[Install]
WantedBy=multi-user.target
